# AAExecutorDaemon/layout/Library/LaunchDaemons/de.fau.cs.dios.pilot.aaexecutord.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamPriority</key>
		<integer>-49</integer>
		<key>JetsamMemoryLimit</key>
		<integer>-1</integer>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>Label</key>
	<string>de.fau.cs.dios.pilot.aaexecutord</string>
	<key>Disabled</key>
	<false/>
	<key>OnDemand</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
	<key>KeepAlive</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/aaexecutord</string>
	</array>
</dict>
</plist>
